# Server/systemctl/monday_discovery.timer
[Unit]
Description=Start Monday discovery every 10 seconds
[Timer]
OnBootSec=10
OnUnitActiveSec=10
AccuracySec=1ms
[Install]
WantedBy=timers.target
